# LaunchAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Install to: ~/Library/LaunchAgents/com.cmuav.maix.kiosk.plist
  Load:       launchctl load -w ~/Library/LaunchAgents/com.cmuav.maix.kiosk.plist
  Replace /Applications/MaixKiosk.app with the final install path.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cmuav.maix.kiosk</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/MaixKiosk.app/Contents/MacOS/MaixKiosk</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    <key>StandardOutPath</key>
    <string>/tmp/maix-kiosk.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/maix-kiosk.err.log</string>
</dict>
</plist>
